# 99-udev-backer.rules.in
KERNEL=="sd*", ATTRS{model}=="@MODEL@", ATTRS{serial}=="@SERIAL@", RUN+="/usr/local/bin/duplicity-backup-runner"
